# profile.d/executable_20-lang-python.sh
#! /usr/bin/env bash
# Python dev utilities

# [uv](https://docs.astral.sh/uv/)
# Install
# brew install uv
# Path settings already set in 10-env-bin.sh
