# profile.d/executable_20-lang-go.sh
#! /usr/bin/env bash
# [Go](https://go.dev/)
if [[ -d $HOME/.gvm ]]; then
    source "$HOME/.gvm/scripts/gvm"
fi
